# community/awesome/desktop-overlay/etc/skel/.config/autostart/picom-detect-once.desktop
[Desktop Entry]
Type=Application
Name=Picom detect once
Exec=/usr/local/bin/picom-detect-once
X-GNOME-Autostart-enabled=true
NoDisplay=true
